# deploy/systemd/sysplane-proxy.service
[Unit]
Description=sysplane-proxy — Sysplane IDC aggregation proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sysplane-proxy --config /etc/sysplane/proxy.yaml
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sysplane-proxy

NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true

[Install]
WantedBy=multi-user.target
