# deploy/systemd/sysplane-agent.service
[Unit]
Description=sysplane-agent — Sysplane physical machine resource agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/sysplane-agent --config /etc/sysplane/agent.yaml
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sysplane-agent

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadOnlyPaths=/

[Install]
WantedBy=multi-user.target
